# prov-shit/ansible/roles/swarm/pack/mesos/base/templates/mesos_cleanup.service
[Service]
Type=simple
ExecStart=rm -rf {{mesos_work_dir}}
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
