# ansible/roles/base/firebird/templates/firebird.service.j2
[Service]
WorkingDirectory={{ firebird_install_dir }}
ExecStart=/bin/bash ./run_firebird.sh
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=firebird
Environment=PHOENIX_PUBLIC_KEY={{ public_keys_dest_dir }}/public_key.pem
Environment=STRIPE_PUBLISHABLE_KEY={{ stripe_publishable_key }}

[Install]
WantedBy=multi-user.target
