# prov-shit/ansible/roles/swarm/pack/vendor/mesos/base/templates/mesos_cleanup.service
[Service]
Type=simple
ExecStart={{usr_local_bin}}/run_mesos_cleanup.sh
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
